# ramdisk/etc/rc.local
#!/bin/sh



startx /usr/bin/qemu-system-x86_64 --hda /dev/mmcblk1 -m 1532 -full-screen -portrait &







exit 0